# iot_driver_linux/systemd/akko-bpf-battery.service
[Unit]
Description=Akko/MonsGeek Keyboard Battery BPF Loader
Documentation=https://github.com/user/monsgeek-m1-v5-tmr
After=sys-subsystem-hid-devices.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/akko-loader --strategy ondemand
ExecStop=/usr/local/bin/akko-loader unload

# BPF loading requires root-equivalent capabilities
# The loader exits after pinning BPF links to /sys/fs/bpf/akko/

# Security hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
